# systemd/psa.timer
[Unit]
Description=Run psa daily

[Timer]
OnCalendar=*-*-* 12:00:00
Persistent=true

[Install]
WantedBy=timers.target
